# rabbitmq.conf
[program:rabbitmq]
command=sh /home/rabbitmq/run_rabbitmq.sh
